# Packer/peon/config/docker-prune.timer
[Unit]
Description=Daily Docker Prune Timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target